# home/.config/systemd/user/check-updates.service
[Unit]
Description=Check for system updates

[Service]
Type=oneshot
ExecStart=/bin/dmd/check-updates --notify
